# deploy/st-card-builder-api.service
[Unit]
Description=ST Card Builder API (Node)
After=network.target st-card-builder-couch.service
Wants=st-card-builder-couch.service

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__DEPLOY_HOME__/st-card-builder/server
Environment=NODE_ENV=production
EnvironmentFile=__DEPLOY_HOME__/st-card-builder/server/.env
ExecStart=/usr/bin/node src/index.mjs
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
